# deploy/lichess/grug-lichess-bot.service
[Unit]
Description=Grug Lichess bot
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=grugbench
Group=grugbench
RemainAfterExit=yes
WorkingDirectory=/srv/grug-lichess/compose
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
